# target/default/target_skeleton/usr/sbin/cron.daily
#!/bin/sh

# Will run every day (at 00:00)
